# scripts/cron.d/jeromelu
# Jeromelu cron schedule — checked in at scripts/cron.d/jeromelu.
# Installed to /etc/cron.d/jeromelu by scripts/lightsail-deploy.sh.
# Edit the source in this repo, redeploy. Do NOT hand-edit /etc/cron.d/jeromelu.
#
# All jobs run as `ubuntu` — the deploy user, member of the docker group,
# owner of /opt/jeromelu and ~/.aws/credentials. Confirmed via
# docs/operations/aws-resource-inventory.md (LIGHTSAIL_USER=ubuntu).
#
# Timezone: the box is UTC. AEST = UTC+10, AEDT = UTC+11. Schedules below
# target AEST clock times — during DST (AEDT) each job fires an hour later
# in local time. Accepted drift; the alternative (CRON_TZ=Australia/Sydney)
# is supported by Debian's cron but skipped for consistency with the
# existing pg-backup convention.
#
# Note on weekly entries: 09:00 AEST = 23:00 UTC the *previous* day. So
# "Monday 09:15 AEST" is Sunday 23:15 UTC, hence dow=0.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Daily channel-stats refresh — 23:00 UTC = 09:00 AEST / 10:00 AEDT.
# ~3 quota units. Snapshots subscriber/video/view counts into channel_metrics.
0 23 * * *      ubuntu  /opt/jeromelu/scripts/scout-refresh.sh channel-stats

# Daily video refresh — 23:15 UTC = 09:15 AEST / 10:15 AEDT.
# ~750 quota units (~180 channels × ~200 videos). Enumerates new videos per
# channel + snapshots video_metrics. Offset 15 min from channel-stats so the
# two don't compete for DB connections. Daily so velocity / breakout detection
# has 1-day resolution; comfortably under the 10k daily quota even at peak.
15 23 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh videos

# Nightly Postgres backup — 16:30 UTC = 02:30 AEST / 03:30 AEDT.
# Streams pg_dump from the postgres container to s3://jeromelu-public-assets/
# backups/postgres/. Retention via S3 lifecycle (14d). The wrapper writes its
# own status line to stdout; cron redirects that into the log file.
30 16 * * *     ubuntu  /opt/jeromelu/scripts/pg-backup.sh >> /var/log/jeromelu/pg-backup.log 2>&1

# Daily cron-health digest — 00:30 UTC = 10:30 AEST / 11:30 AEDT.
# Runs *after* all the above so it can report on the trailing 24h. Emails
# a per-job status table (run? exit clean? what did it write?) via SES.
# Sister to .github/workflows/cost-report.yml — different signal, same shape.
30 0 * * *      ubuntu  /opt/jeromelu/scripts/cron-report.sh >> /var/log/jeromelu/cron-report.log 2>&1
