# scripts/cron.d/jeromelu
# Jeromelu cron schedule — checked in at scripts/cron.d/jeromelu.
# Installed to /etc/cron.d/jeromelu by scripts/lightsail-deploy.sh.
# Edit the source in this repo, redeploy. Do NOT hand-edit /etc/cron.d/jeromelu.
#
# All jobs run as `ubuntu` — the deploy user, member of the docker group,
# owner of /opt/jeromelu and ~/.aws/credentials. Confirmed via
# docs/operations/aws-resource-inventory.md (LIGHTSAIL_USER=ubuntu).
#
# Timezone: the box is UTC. AEST = UTC+10, AEDT = UTC+11. Schedules below
# target AEST clock times — during DST (AEDT) each job fires an hour later
# in local time. Accepted drift; the alternative (CRON_TZ=Australia/Sydney)
# is supported by Debian's cron but skipped for consistency with the
# existing pg-backup convention.
#
# Note on weekly entries: 09:00 AEST = 23:00 UTC the *previous* day. So
# "Monday 09:15 AEST" is Sunday 23:15 UTC, hence dow=0.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Daily channel-stats refresh — 23:00 UTC = 09:00 AEST / 10:00 AEDT.
# ~3 quota units. Snapshots subscriber/video/view counts into channel_metrics.
0 23 * * *      ubuntu  /opt/jeromelu/scripts/scout-refresh.sh channel-stats

# Daily video refresh — 23:15 UTC = 09:15 AEST / 10:15 AEDT.
# ~750 quota units (~180 channels × ~200 videos). Enumerates new videos per
# channel + snapshots video_metrics (change-only: records a row only when a
# video's views/likes/comments changed since the last snapshot; ~70% of daily
# samples are unchanged and skipped). Offset 15 min from channel-stats so the
# two don't compete for DB connections. Daily so velocity / breakout detection
# has 1-day resolution; comfortably under the 10k daily quota even at peak.
15 23 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh videos

# Weekly deterministic YouTube source discovery — Sunday 20:30 UTC = Monday
# 06:30 AEST / 07:30 AEDT. Uses the default deterministic query bank with
# max_results=10 and max_videos=25; costs roughly 1.3k YouTube quota units
# (13 search.list calls plus low-cost metadata batches). Weekly keeps the
# candidate inbox fresh without crowding the daily YouTube refresh budget.
30 20 * * 0     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh source-discovery-youtube

# Daily SuperCoach roster refresh — 22:45 UTC = 08:45 AEST / 09:45 AEDT.
# Fetches the canonical players-cf roster, archives to S3, and applies the
# SCD-2 people/player_attributes refresh via /api/admin/scout/supercoach-roster.
45 22 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh supercoach-roster

# Current-round SuperCoach stats — Sunday/Tuesday/Thursday 22:55 UTC =
# Monday/Wednesday/Friday 08:55 AEST. Resolves current_round from the public
# SuperCoach settings JSON, then calls /api/admin/scout/supercoach-stats.
# Idempotent UPSERT on player_rounds; no historical backfill here.
55 22 * * 0,2,4 ubuntu  /opt/jeromelu/scripts/scout-refresh.sh supercoach-stats current

# Weekly SuperCoach teams refresh — Mondays 23:30 UTC = Tuesday 09:30 AEST.
# Tiny payload (17 rows, ~3KB). Refreshes teams.metadata_json.supercoach
# via /api/admin/scout/supercoach-teams. Audit row under agent_id='scout'.
30 23 * * 1     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh supercoach-teams

# Weekly SuperCoach settings snapshot — Mondays 23:35 UTC = Tuesday 09:35 AEST.
# ~15KB payload. Captures game rules per season into sc_settings (classic
# mode only — draft is on-demand via make scout-supercoach-settings MODE=draft).
35 23 * * 1     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh supercoach-settings

# Weekly nrl.com players roster — Mondays 23:40 UTC = Tuesday 09:40 AEST.
# Walks the 17 NRL teams via refresh-all (1 req/sec; ~20s wall time). Archives
# scout/nrlcom/players-roster/111/team-{team_id}.json per team. Profile data
# (DOB, position, image) changes slowly; weekly is right.
40 23 * * 1     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-players-roster

# Daily nrl.com draw refresh — 18:00 UTC = 04:00 AEST. Current-round fixtures for
# NRL (111), current season. Archives scout/nrlcom/draw/111/{season}/round-NN.json.
0 18 * * *      ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-draw

# Daily nrl.com match-centre walk — 18:15 UTC = 04:15 AEST. Walks the current
# round's fixtures, archives each match's full JSON. Round resolved server-side.
15 18 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-match-centre

# Daily nrl.com casualty-ward snapshot — 18:30 UTC = 04:30 AEST. Captures the
# official league injury roll for NRL (111). Timestamped S3 key per day
# (scout/nrlcom/casualty-ward/111/{YYYYMMDD}.json) — preserves the daily
# timeline the injuries state-machine extractor replays. No season param.
30 18 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-casualty-ward

# Daily nrl.com ladder snapshot — 18:45 UTC = 04:45 AEST. Current-round team
# standings + the 22 per-team metrics for NRL (111), current season. Archives
# scout/nrlcom/ladder/111/{season}/round-NN.json (round resolved server-side).
45 18 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-ladder

# Daily nrl.com stats snapshot — 18:50 UTC = 04:50 AEST. Pre-computed top-25
# leaderboards (~275KB) for NRL (111), current season. Archives
# scout/nrlcom/stats/111/{season}.json. Cheap to run daily even though the
# leaderboards only change a few times per week.
50 18 * * *     ubuntu  /opt/jeromelu/scripts/scout-refresh.sh nrlcom-stats

# Daily Scout S3-to-DB projection — 19:20 UTC = 05:20 AEST. Runs after the
# nrl.com archive jobs above and projects the current-season S3 captures into
# rounds, matches, team lists, player match stats, timelines, standings,
# leaderboards, injuries, and derived nrl.com attributes. Wrapper logs to
# /var/log/jeromelu/scout-populate.log and stages scripts into jeromelu-api.
20 19 * * *     ubuntu  /opt/jeromelu/scripts/scout-populate.sh nrlcom-current

# Nightly Postgres backup — 16:30 UTC = 02:30 AEST / 03:30 AEDT.
# Streams pg_dump from the postgres container to s3://jeromelu-public-assets/
# backups/postgres/. Retention via S3 lifecycle (14d). The wrapper writes its
# own status line to stdout; cron redirects that into the log file.
30 16 * * *     ubuntu  /opt/jeromelu/scripts/pg-backup.sh >> /var/log/jeromelu/pg-backup.log 2>&1

# Daily cron-health digest — 00:30 UTC = 10:30 AEST / 11:30 AEDT.
# Runs *after* all the above so it can report on the trailing 24h. Emails
# a per-job status table (run? exit clean? what did it write?) via SES.
# Sister to .github/workflows/cost-report.yml — different signal, same shape.
30 0 * * *      ubuntu  /opt/jeromelu/scripts/cron-report.sh >> /var/log/jeromelu/cron-report.log 2>&1

# Daily app-error digest — 00:35 UTC = 10:35 AEST / 11:35 AEDT.
# Scrapes `docker logs --since 24h` for api + web containers, counts
# error lines + exception classes, surfaces samples. Lands five min
# after cron-report so the morning batch arrives together.
35 0 * * *      ubuntu  /opt/jeromelu/scripts/error-report.sh >> /var/log/jeromelu/error-report.log 2>&1

# Weekly content digest — Mondays 22:00 UTC = Tuesday 08:00 AEST.
# New videos, channel velocity, claims-by-type, predictions resolved.
# Different signal from the ops emails — this is the product-owner view.
0 22 * * 1      ubuntu  /opt/jeromelu/scripts/content-report.sh >> /var/log/jeromelu/content-report.log 2>&1

# Weekly capacity report — Mondays 22:30 UTC = Tuesday 08:30 AEST.
# `df -h /`, `pg_database_size`, top tables by size. Catches "you're
# 60 days from full" before a crash does.
30 22 * * 1     ubuntu  /opt/jeromelu/scripts/disk-report.sh >> /var/log/jeromelu/disk-report.log 2>&1
